# deploy/makearoll-deploy.timer
[Unit]
Description=Run makearoll_discord redeploy check every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
